# templates/spark-worker.conf.erb
description "Spark Worker Service script"
start on runlevel [2345]
stop on runlevel [06]
setuid root
setgid root
console log
chdir <%= @install_dir %>



exec <%= @install_dir %>/bin/spark-worker-runner.sh

