# templates/spark-master.conf.erb
description "Spark Master Service script"
start on runlevel [2345]
stop on runlevel [06]
setuid spark
setgid spark
console log
chdir <%= @install_dir %>

exec <%= @install_dir %>/bin/spark-master-runner.sh

